# root/etc/services.d/guacd/run
#!/usr/bin/with-contenv sh

GUACD_LOGLEVEL=debug

echo "Starting guacamole guacd..."
s6-setuidgid root guacd -b 0.0.0.0 -l 4822 -L $GUACD_LOGLEVEL -f
